# lib/supervisor.conf
[supervisord]
nodaemon=true

[program:mongod]
command=/usr/bin/mongo

[program:apache2]
command=/bin/bash -c "source /etc/apache2/envvars && exec /usr/sbin/apache2 -DFOREGROUND"